# deploy/kleats.service
[Unit]
Description=KLEats Next.js App
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/kleatsv2
EnvironmentFile=/root/kleatsv2/.env.local
Environment=NODE_ENV=production
Environment=FRONTEND_PORT=3001
ExecStart=/usr/bin/node /root/kleatsv2/server.mjs
Restart=always
RestartSec=10
User=root
Group=root

[Install]
WantedBy=multi-user.target
